# prowlarr/rootfs/etc/s6-overlay/s6-rc.d/prowlarr/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Prowlarr
# Runs Prowlarr
# ==============================================================================

bashio::log.info 'Starting Prowlarr...'

exec /opt/bin/Prowlarr -nobrowser -data=/config
